# rpmbuild/SOURCES/pure-ftpd
#!/bin/bash
# chkconfig: 345 99 01
# description: pure-ftpd init script

. /etc/profile

PIDFILE=/var/run/pure-ftpd.pid

case $1 in
start)

/usr/sbin/pure-ftpd &
echo "pure-ftpd started"
;;
stop)

echo "pure-ftpd stopped"	
kill $(cat "$PIDFILE") && rm -f "$PIDFILE"
;;
restart)
$0 stop
sleep 5;
$0 start
;;
*)
echo "Usage $0 start|stop|restart"
esac
exit 0

